# infra/voicedrop-en-sync/com.jianshuo.voicedrop-en-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- voicedrop 中→英夜间同步 launchd 模板。安装（需人工执行，因为 sync.sh 内含
     无人值守的 claude 调用，权限门要求本人确认）：
       cp infra/voicedrop-en-sync/com.jianshuo.voicedrop-en-sync.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.jianshuo.voicedrop-en-sync.plist
     卸载：launchctl unload ~/Library/LaunchAgents/com.jianshuo.voicedrop-en-sync.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jianshuo.voicedrop-en-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jianshuo/code/jianshuo.dev/infra/voicedrop-en-sync/sync.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>30</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/jianshuo/.claude/automation/voicedrop-en-sync.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jianshuo/.claude/automation/voicedrop-en-sync.stderr.log</string>
</dict>
</plist>
