# claude-agent/deploy/claude-agent.service
[Unit]
Description=Claude Agent SDK web chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-agent
Group=claude-agent
WorkingDirectory=/opt/claude-agent
Environment=NODE_ENV=production
Environment=HOME=/opt/claude-agent
EnvironmentFile=/opt/claude-agent/.env
ExecStart=/usr/bin/node /opt/claude-agent/dist/server.js
Restart=always
RestartSec=3

# --- hardening: confine the agent's bash/file tools' blast radius ---
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/claude-agent
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
