# codex-agent/deploy/codex-agent.service
[Unit]
Description=codex-agent - Codex subscription web chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codex-agent
Group=codex-agent
WorkingDirectory=/opt/codex-agent
Environment=NODE_ENV=production
Environment=HOME=/opt/codex-agent
# 独立登录会话（2026-07-05 修订：paint 的 gpt-image-2-skill 与 codex CLI 对
# auth.json 的 last_refresh 字段格式互不兼容——epoch vs RFC3339，不能共享文件。
# 独立 codex login 有自己的 refresh token 链，与 paint 互不踢）。
Environment=CODEX_HOME=/opt/codex-agent/.codex
EnvironmentFile=-/opt/codex-agent/.env
ExecStart=/usr/bin/node /opt/codex-agent/dist/server.js
Restart=always
RestartSec=3

# B 档权限：允许 sudo 白名单（所以不能开 NoNewPrivileges/RestrictSUIDSGID），
# 写权限限定在自身目录、共享凭证、/etc/caddy（管 Caddy 配置用）
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/codex-agent /etc/caddy
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
